# server/pharos-brief.timer
[Unit]
Description=Trigger Pharos daily brief at 15:00 Asia/Taipei

[Timer]
# 15:00 Taipei time = 07:00 UTC
OnCalendar=*-*-* 07:00:00 UTC
Persistent=true
Unit=pharos-brief.service

[Install]
WantedBy=timers.target
