# server/yt2epub-consolidate-vocab.timer
[Unit]
Description=Weekly vocabulary consolidation fallback
Requires=yt2epub-consolidate-vocab.service

[Timer]
# 週一早上 11:00（Mac 中午跑之前 1 小時、若 Mac 上週都沒跑就先兜底）
OnCalendar=Mon *-*-* 11:00:00
Persistent=true

[Install]
WantedBy=timers.target
